# src/build/api/support-files/templates/#etc#supervisor-usermgr-api.conf
[unix_http_server]
file=__BK_HOME__/logs/usermgr/supervisor.sock

[supervisord]
logfile=__BK_HOME__/logs/usermgr/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=__BK_HOME__/logs/usermgr/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://__BK_HOME__/logs/usermgr/supervisor.sock

[program:usermgrapi]
directory=__BK_HOME__/usermgr/api
environment=
    PATH="__BK_HOME__/.envs/usermgr-api/bin/",
    DJANGO_SETTINGS_MODULE="bkuser_core.config.overlays.prod",
    LOGS_PATH="__BK_HOME__/logs",
command=__BK_HOME__/.envs/usermgr-api/bin/gunicorn wsgi -w 16 --max-requests 1024 --max-requests-jitter 50 --threads 2 --worker-class gevent -b :8009 --access-logfile - --error-logfile -
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/usermgr/stdout.log


[program:usermgrapi-worker]
directory=__BK_HOME__/usermgr/api
environment=
    PATH="__BK_HOME__/.envs/usermgr-api/bin/",
    DJANGO_SETTINGS_MODULE="bkuser_core.config.overlays.prod",
    LOGS_PATH="__BK_HOME__/logs",
command=__BK_HOME__/.envs/usermgr-api/bin/celery -A bkuser_core worker -l info --concurrency=8 --max-tasks-per-child=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/usermgr/worker.log


[program:usermgrapi-beat]
directory=__BK_HOME__/usermgr/api
environment=
    PATH="__BK_HOME__/.envs/usermgr-api/bin/",
    DJANGO_SETTINGS_MODULE="bkuser_core.config.overlays.prod",
    LOGS_PATH="__BK_HOME__/logs",
command=__BK_HOME__/.envs/usermgr-api/bin/celery -A bkuser_core beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=__BK_HOME__/logs/usermgr/beat.log
